# docker/supervisor-app.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /home/g3wsuite/uwsgi.ini
stopsignal = QUIT

[program:nginx-app]
command = /usr/sbin/nginx